# packaging/alarm-session-agent.socket
[Socket]
ListenStream=/run/alarm_agent/%U
DirectoryMode=0777
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@

[Install]
WantedBy=sockets.target
